# rpi4_print_worker.service
[Unit]
Description=Restaurant Print Queue Worker
After=network.target
Wants=network.target

[Service]
Type=simple
User=gsz
Group=gsz
WorkingDirectory=/home/gsz
ExecStart=/usr/bin/python3 /home/gsz/rpi4_print_queue_worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Variables de entorno
Environment=PYTHONPATH=/home/gsz
Environment=API_BASE_URL=http://restaurant-backend:8000/api

[Install]
WantedBy=multi-user.target